# PROJECTS/whatsapp-agentic-bot/com.hotelbot.whatsapp-agentic-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hotelbot.whatsapp-agentic-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/npm</string>
        <string>run</string>
        <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/g/Desktop/MAIA opencode/whatsapp-agentic-bot</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/g/Desktop/MAIA opencode/whatsapp-agentic-bot/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/g/Desktop/MAIA opencode/whatsapp-agentic-bot/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>

    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>10240</integer>
    </dict>
</dict>
</plist>
